# systemd/user.control/sad.service
[Unit]
Description=SAD music player daemon

[Service]
ExecStart=/usr/local/bin/sad
Restart=on-failure

[Install]
WantedBy=basic.target
